# conf/supervisor.conf
[program:dp]
command=python app.py
numprocs=1
autorestart=true
autostart=true
directory=/var/www/adjust_deeplink/
user=www-data
stderr_logfile=/var/log/dp.err.log
stdout_logfile=/var/log/dp.out.log